# roles/apps/timeless/templates/timeless-docker.upstart
description     "Timeless master"
author          "Gabriel Falcao <gabriel@nacaolivre.org>"

start on (started networking)

respawn
respawn limit 5 60

env CASSANDRA_HOST={{ timeless_cassandra_host }}
env CASSANDRA_PORT={{ timeless_cassandra_port }}
env CASSANDRA_PORT={{ timeless_cassandra_port }}

exec docker run --rm --net=host \
  -e CASSANDRA_PORT -e CASSANDRA_HOST \
  --name=timeless \
  --volume=/srv/timeless:/srv/timeless \
  timeless
