# roles/quietness/templates/cassandra.upstart
description     "Cassandra"
author          "Gabriel Falcao <gabriel@nacaolivre.org>"

start on (started networking)

respawn
respawn limit 5 60


exec docker run --rm --net=host \
  --log-driver="json-file" \
  --volume=/srv/cassandra:/data \
  cassandra
