# roles/quietness/templates/nginx-quietness-docker.upstart
description     "Nginx quietness"
author          "Gabriel Falcao <gabriel@nacaolivre.org>"

start on (started networking)

respawn
respawn limit 5 60

exec docker run --rm --net=host \
  -e CASSANDRA_PORT -e CASSANDRA_HOST \
  --volume=/srv/quietness:/srv/quietness \
  nginx
