# roles/quietness/templates/quietness-docker.upstart
description     "Quietness master"
author          "Gabriel Falcao <gabriel@nacaolivre.org>"

start on (started networking)

respawn
respawn limit 5 60

env CASSANDRA_HOST={{ quietness_cassandra_host }}
env CASSANDRA_PORT={{ quietness_cassandra_port }}

exec docker run --rm --net=host \
  -e CASSANDRA_PORT -e CASSANDRA_HOST \
  --name=quietnessapi \
  --volume=/srv/quietness:/srv/quietness \
  quietnessapi

EXPOSE 5000